# etc/systemd/system/named.service
[Unit]
Description=Named

[Service]
Type=oneshot
ExecStart=/usr/sbin/named

[Install]
WantedBy=multi-user.target